# mysql-backup/crontab
SHELL=/bin/bash
0 * * * * /scripts/run.sh
0 0 * * 0 logrotate /etc/logrotate.d/mysql-backup
